# linux/systemd/crossring.service
[Unit]
Description=CROSSRING Endpoint Security Daemon
After=network.target apparmor.service

[Service]
Type=simple
ExecStart=/usr/bin/crossring-daemon
Restart=always
RestartSec=5
User=root
Group=root

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/crossring /run
PrivateTmp=true
NoNewPrivileges=false
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_NET_ADMIN CAP_SYS_ADMIN
AmbientCapabilities=CAP_SYS_PTRACE CAP_NET_ADMIN

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crossring

[Install]
WantedBy=multi-user.target
